# client/usr/lib/systemd/system/prometheus-tunnel.service.d/tunnel.conf
[Service]
PROMETHEUS_TUNNEL_USER=prometheus-tunnel
PROMETHEUS_TUNNEL_HOST=prometheus.example.com
PROMETHEUS_TUNNEL_PORT=22
