# client/usr/lib/systemd/system/prometheus-tunnel.service.d/tunnel.conf
[Service]
Environment=PROMETHEUS_TUNNEL_USER=prometheus-tunnel
Environment=PROMETHEUS_TUNNEL_HOST=prometheus.example.com
Environment=PROMETHEUS_TUNNEL_PORT=22
